# packages/pipa-fixes/pkg/pipa-fixes/usr/lib/systemd/system/pipa-bt-mac.service
[Unit]
Description=Set BT mac address
Wants=bluetooth.target
After=bluetooth.target

[Service]
Type=oneshot
ExecStart=/usr/bin/pipa-bt-mac-set
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
