# packages/pipa-fixes/pkg/pipa-fixes/usr/lib/systemd/system/pipa-wlan-mac.service
[Unit]
Description=Set WLAN mac address based on board serial
Wants=network-pre.target
Before=network-pre.target
BindsTo=sys-subsystem-net-devices-wlp1s0.device
After=sys-subsystem-net-devices-wlp1s0.device

[Service]
Type=oneshot
ExecStart=/usr/bin/pipa-wlan-mac-set
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
